# common/morpher-librato.cron
# Send metrics to Librato
* * * * * morpher /usr/bin/morpher-librato &> /dev/null || :
