# common/morpher-librato.cron
SHELL=/bin/bash
HOME=/srv/projects/pkgre

* * * * * morpher /usr/bin/morpher-librato &> /dev/null || :
